# scripts/do_phase1/systemd/clypt-phase1-visual.service
[Unit]
Description=Clypt Phase 1 Visual Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/clypt-phase1/repo
EnvironmentFile=/etc/clypt-phase1/phase1.env
EnvironmentFile=/etc/clypt/clypt-phase1-runtime.env
Environment=PYTORCH_KERNEL_CACHE_PATH=/opt/clypt-phase1/.cache/torch/kernels
ExecStart=/opt/clypt-phase1/venvs/phase1/bin/python -m backend.runtime.run_phase1_visual_service --host 127.0.0.1 --port 9200
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
