# scripts/do_phase1/systemd/clypt-vllm-qwen.service
[Unit]
Description=Clypt local Qwen vLLM service
After=docker.service nvidia-persistenced.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=5
TimeoutStartSec=3600
ExecStartPre=-/usr/bin/docker stop clypt-vllm-qwen
ExecStartPre=-/usr/bin/docker rm clypt-vllm-qwen
ExecStart=/usr/bin/docker run --name clypt-vllm-qwen --gpus all --ipc=host -p 127.0.0.1:8001:8000 -v /opt/clypt-phase1/hf-cache:/root/.cache/huggingface vllm/vllm-openai:latest --model Qwen/Qwen3.5-27B --served-model-name Qwen/Qwen3.5-27B --language-model-only --gpu-memory-utilization 0.55 --max-num-seqs 12 --max-num-batched-tokens 12288 --max-model-len 131072 --kv-cache-dtype fp8_e4m3 --enable-prefix-caching
ExecStop=/usr/bin/docker stop clypt-vllm-qwen

[Install]
WantedBy=multi-user.target
