# scripts/do_phase1/systemd/clypt-phase1-vllm-vibevoice.service
[Unit]
Description=Clypt Phase 1 VibeVoice vLLM Sidecar
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/clypt-phase1/repo
EnvironmentFile=/etc/clypt-phase1/phase1.env
ExecStartPre=-/usr/bin/docker stop clypt-phase1-vllm-vibevoice
ExecStartPre=-/usr/bin/docker rm clypt-phase1-vllm-vibevoice
ExecStart=/usr/local/bin/clypt-phase1-run-vllm-vibevoice
ExecStop=/usr/bin/docker stop clypt-phase1-vllm-vibevoice
Restart=always
RestartSec=10
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
