# scripts/do_phase1/systemd/clypt-phase1-worker.service
[Unit]
Description=Clypt Phase 1 DigitalOcean Worker
After=network-online.target clypt-phase1-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/clypt-phase1/repo
ExecStart=/opt/clypt-phase1/repo/scripts/do_phase1/run_worker_service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
