# scripts/do_phase1_visual/systemd/clypt-v31-phase24-local-worker.service
[Unit]
Description=Clypt V3.1 Phase24 local worker
After=network.target clypt-sglang-qwen.service
Requires=clypt-sglang-qwen.service
RequiresMountsFor=/opt/clypt-phase1/repo

[Service]
Type=simple
WorkingDirectory=/opt/clypt-phase1/repo
EnvironmentFile=/etc/clypt-phase1/v3_1_phase1.env
Environment=HOME=/opt/clypt-phase1
Environment=CLYPT_PHASE1_CACHE_HOME=/opt/clypt-phase1/.cache
Environment=XDG_CACHE_HOME=/opt/clypt-phase1/.cache
Environment=TORCH_HOME=/opt/clypt-phase1/.cache/torch
Environment=HF_HOME=/opt/clypt-phase1/.cache/huggingface
ExecStart=/opt/clypt-phase1/venvs/phase1/bin/python -m backend.runtime.run_phase24_local_worker
Restart=always
RestartSec=5
LimitNOFILE=524288
CPUAccounting=yes
MemoryAccounting=yes

[Install]
WantedBy=multi-user.target
