# scripts/do_phase1_audio/systemd/clypt-audio-host.service
[Unit]
Description=Clypt VibeVoice ASR + node-media-prep service (RTX 6000 Ada)
After=network-online.target clypt-vllm-vibevoice.service
Wants=network-online.target
Requires=clypt-vllm-vibevoice.service

[Service]
Type=simple
WorkingDirectory=/opt/clypt-audio-host/repo
EnvironmentFile=/etc/clypt-audio-host/audio_host.env
Environment=HOME=/opt/clypt-audio-host
Environment=CLYPT_PHASE1_CACHE_HOME=/opt/clypt-audio-host/.cache
Environment=XDG_CACHE_HOME=/opt/clypt-audio-host/.cache
Environment=CLYPT_PHASE1_VIBEVOICE_ASR_SERVICE_SCRATCH_ROOT=/opt/clypt-audio-host/scratch
ExecStart=/opt/clypt-audio-host/venvs/audio/bin/python -m uvicorn backend.runtime.phase1_audio_service.app:create_app --factory --host ${CLYPT_PHASE1_VIBEVOICE_ASR_SERVICE_BIND} --port ${CLYPT_PHASE1_VIBEVOICE_ASR_SERVICE_PORT} --workers 1 --timeout-keep-alive 120
Restart=always
RestartSec=5
LimitNOFILE=524288
CPUAccounting=yes
MemoryAccounting=yes

[Install]
WantedBy=multi-user.target
